# launchd/com.agent-pager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agent-pager</string>
  <key>ProgramArguments</key>
  <array>
    <!-- These paths are replaced by setup.sh -->
    <string>/usr/local/bin/node</string>
    <string>/Users/piyush/GitHub/claude-slack-bridge/bridge.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/piyush/GitHub/claude-slack-bridge</string>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/tmp/agent-pager.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/agent-pager.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
